# scripts/macosx/init/profile.d/ssh.sh
#!/bin/bash
if [ ! -d "/tmp/ssh" ]; then
   mkdir /tmp/ssh
fi

