# scripts/os/macosx/init/profile.d/go.sh
export GOROOT=/usr/local/opt/go/libexec
export GOPATH=$HOME/go-projects
export PATH=$PATH:$GOROOT/bin:$HOME/go-projects/bin
