# scripts/os/macosx/init/profile.d/java.sh
export JAVA_HOME=$HOME/opt/java
export MVN_HOME=$HOME/opt/mvn
export MAVEN_OPTS="-Xms512m -Xmx1024m"
export PATH=$MVN_HOME/bin:$PATH
