# scripts/os/macosx/init/profile.d/libvirt.sh
export LIBVIRT_AUTH_FILE=~/.libvirt/auth.conf
