# init/profile.d/git.sh
#git duet
export GIT_DUET_SECONDS_AGO_STALE=30
